# target/init_scripts/actlogo.service
[Unit] 
Description=Terminate Psplash Act Screen
After=bootlogo.service

[Service] 
Type=forking 
ExecStart=/bin/bash /opt/scripts/boot/psplash_quit.sh

[Install]
WantedBy=multi-user.target