# target/init_scripts/bootlogo.service
[Unit]
Description=Starts Psplash Boot logo
DefaultDependencies=no
After=systemd-udev-trigger.service systemd-udevd.service

[Service]
ExecStart=/bin/bash /opt/scripts/boot/psplash.sh

[Install]
WantedBy=sysinit.target
