# target/init_scripts/solve_qt_deb.service
[Unit]
Description=solve QT dependence

[Service]
Type=oneshot
ExecStart=/bin/sh /opt/scripts/boot/ebf-build.sh

[Install]
WantedBy=multi-user.target