# target/init_scripts/systemd-generic-board-startup.service
[Unit]
Description=Generic Board Startup
After=local-fs.target solve_qt_deb.service 
Before=serial-getty@.service
[Service]
Type=oneshot
ExecStart=/bin/sh /opt/scripts/boot/generic-startup.sh

[Install]
WantedBy=multi-user.target
