# target/init_scripts/systemd-generic-board-startup.service
[Unit]
Description=Generic Board Startup
After=local-fs.target 
Before=serial-getty@.service
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/sh /opt/scripts/boot/generic-startup.sh
ExecStop=/sbin/modprobe -r g_multi
[Install]
WantedBy=multi-user.target
